# couture/resources/couture_emails.service
[Unit]
Description=Couture hourly

[Service]
WorkingDirectory=/srv/couture/
ExecStart=/usr/bin/nodejs /srv/couture/scripts/order_emails.js
